# vm/root/guestproxy.service
[Unit]
Description=Guest Proxy

[Service]
ExecStart=/usr/local/bin/guest

[Install]
WantedBy=multi-user.target
